# installer/etc/init/copy-disk-image.conf

description "disk image installation task"

start on runlevel [2345]
task
exec /etc/init/copy-disk-image.sh

