# installer/etc/init/copy-disk-image.conf

description "disk image installation task"

start on local-filesystem
task
exec /etc/init/copy-disk-image.sh

